# vagrant-azure/production-webconfig/supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn hearcloud.wsgi --bind 0.0.0.0:8000
directory=/home/hearclouduser/hear-cloud
user=hearclouduser
autostart=true
autorestart=true
redirect_stderr=true
