# vagrant-azure/production-webconfig/supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn hearcloud.wsgiprod --bind 0.0.0.0:8000
directory=/home/vagrant/hear-cloud
user=vagrant
autostart=true
autorestart=true
redirect_stderr=true
